# LarbsP/.config/systemd/user/mbsync.service
[Unit]
Description=Sync mail with mbsync

[Service]
ExecStart=/usr/bin/mbsync -a